# systemd/wiimote-bridge.service
[Unit]
Description=Wiimote to USB Gamepad Bridge
After=bluetooth.target
Wants=bluetooth.target wiimote-gadget.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStartPre=/usr/bin/hciconfig hci0 up piscan
ExecStart=/usr/bin/python3 /opt/wiimote-bridge/wiimote_bridge.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wiimote-bridge

# Run as root (needed for /dev/hidg* access and Bluetooth)
User=root

# Hardening
ProtectHome=read-only
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
